# systemd/motion.service
[Unit]
Description=Person Detection and Recording
After=network-online.target
Wants=network-online.target

[Service]
# --- Your Python environment ---
WorkingDirectory=/home/yourname
ExecStart=/home/yourname/yolo-infer/bin/python /home/yourname/person.py

# --- Restart behavior ---
Restart=on-failure
RestartSec=2

# --- Shutdown behavior ---
TimeoutStopSec=10
KillMode=mixed

# --- Logging ---
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

# --- Security / isolation (optional but recommended) ---
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

